# laurel-api.service
[Unit]
Description=Laurel clip management API
After=network.target

[Service]
Type=simple
ExecStart=/opt/laurel-server/laurel-api
Environment=CLIP_DIR={{CLIP_DIR}}
Environment=BIND={{API_BIND}}
Environment=PORT={{API_PORT}}
User={{USER}}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
